# crontab/analysis.sh
#!/usr/bin/env bash

cd /home/zhaojingming/apps/log_analysis/
sh scripts/analysis.sh >> log/analysis.log